# gpu-coder-frontend.service
[Unit]
Description=GPU coder dashboard
After=network.target

[Service]
Type=simple
User=vishwa
WorkingDirectory=/home/vishwa/gpu/dock-dash-command-center
ExecStart=/bin/bash -c "cd /home/vishwa/gpu/dock-dash-command-center && npm run dev"
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpu-coder-dashboard

[Install]
WantedBy=multi-user.target
