# backend/gpu-coder-admin.service
[Unit]
Description=GPU Coder Backend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=vishwa
Group=vishwa
WorkingDirectory=/home/vishwa/gpu/dock-dash-command-center/backend
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/home/vishwa/gpu/dock-dash-command-center/backend
ExecStart=/bin/bash -c "cd /home/vishwa/gpu/dock-dash-command-center/backend && source qvp-deploy/bin/activate && python app.py"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpu-coder-backend

[Install]
WantedBy=multi-user.target
