# backend/scripts/cxl-manager.service
[Unit]
Description=CXL Manager Services (Auth Service + Web Interface)
After=network.target

[Service]
User=vishwa
Group=vishwa
WorkingDirectory=/home/vishwa/workspace/cxl-deploy/docker-code-server/manager
Environment="PYTHONPATH=/home/vishwa/workspace/cxl-deploy/docker-code-server"
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/home/vishwa/.local/bin:$PATH"
ExecStart=/bin/bash -c 'source /path/to/python/env/bin/activate && python3 auth_service.py & streamlit run web_interface.py'
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cxl-manager

[Install]
WantedBy=multi-user.target