# backend/agent/gpu-coder-agent.service
[Unit]
Description=GPU Coder Agent Service
After=network.target docker.service
Wants=network.target
Requires=docker.service

[Service]
Type=simple
User=vishwa
Group=vishwa
WorkingDirectory=/home/vishwa/gpu/dock-dash-command-center/backend/agent
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/home/vishwa/gpu/dock-dash-command-center/backend/agent
ExecStart=/bin/bash -c "cd /home/vishwa/gpu/dock-dash-command-center/backend && source qvp-deploy/bin/activate && python agent/agent_server.py"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpu-coder-agent

[Install]
WantedBy=multi-user.target
